# systemd/heizungSetRelays.timer
[Unit]
Description=Set relays every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
